# web-nginx-mysql/ubuntu/conf/etc/supervisor/conf.d/supervisord_zabbix.conf
[supervisord]
nodaemon = true

[program:nginx]
command = /usr/sbin/nginx -g "daemon off;" -c /etc/nginx/nginx.conf
auto_start = true
autorestart = true

stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

[program:php-fpm]
command = /usr/sbin/php5-fpm -F -c /etc/php5/fpm/php-fpm.conf
auto_start = true
autorestart = true

stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0